# netguard/install/netguard-monitor.service
[Unit]
Description=NetGuard bandwidth monitor and enforcer
After=graphical-session.target network.target polkit.service
Wants=polkit.service
# pkexec's allow_active=yes check depends on logind having already marked this
# session active; that can race with this service on login. Give it more
# retry headroom than the default before systemd gives up entirely.
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
ExecStart=%h/.local/share/netguard/venv/bin/python -m netguard.monitor.collector
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
